# crontab
*/5 * * * * cd /app && ./locomotive -vvv $(printenv REMOTE_SERVER) >> /root/.config/locomotive-cron.log 2>&1
